# crontab
# Auto-renew SSL certificates daily at 2:30 AM
30 2 * * * curl -X POST http://localhost:3000/api/certificates/renew-all >/dev/null 2>&1 